# packages/piccolo-os-support/piccolo-clock-epoch-save.service
[Unit]
Description=Piccolo clock epoch save

[Service]
Type=oneshot
ExecStart=/usr/libexec/piccolo/clock-epoch.sh save
