# packages/piccolo-os-support/piccolo-watchdog-check.service
[Unit]
Description=Piccolo watchdog driver validation

[Service]
Type=oneshot
ExecStart=/usr/libexec/piccolo/watchdog-check.sh

[Install]
WantedBy=multi-user.target
