# src/l0/kiwi/root/etc/systemd/system/piccolod.service
[Unit]
Description=Piccolo Orchestrator (piccolod)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/piccolo/current/bin/piccolod serve --config /etc/piccolo/config.yaml
Restart=always
RestartSec=2s
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
