# src/l0/kiwi/prod/root/etc/systemd/system/piccolo-rollback.service
[Unit]
Description=Piccolo Automatic Rollback on Health Failure
After=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/piccolo-rollback.sh
TimeoutSec=15min
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/run /var
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service

# No [Install]; triggered by OnFailure from health-check

