# src/l0/kiwi/prod/root/etc/systemd/system/piccolod-health-check.service
[Unit]
Description=Piccolo Health Check for OS Updates - PRODUCTION MODE
Documentation=https://piccolospace.com/docs/health-checks
After=piccolod.service
# Run even if piccolod fails to start so we can trigger automatic rollback.
Wants=piccolod.service
OnFailure=piccolo-rollback.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/piccolod-health-check-prod.sh
TimeoutSec=60
RemainAfterExit=yes

# PRODUCTION SECURITY HARDENING
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadOnlyPaths=/
ReadWritePaths=/run /var/log
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
PrivateDevices=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=yes

# PRODUCTION LOGGING - Minimal output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piccolod-healthcheck-prod

# Allow writing sentinel state
ReadWritePaths=/run /var/log /var/lib/piccolo

[Install]
WantedBy=multi-user.target
